# lf/.config/systemd/user/clear-lf-cache.timer
[Unit]
Description=Triggers a service to clear lf cache

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
